# ops/launchd/com.frikshun.creator-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.frikshun.creator-tunnel</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/ngrok</string>
    <string>http</string>
    <string>5050</string>
    <string>--url</string>
    <string>https://sneer-runny-old.ngrok-free.dev</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/allentaylor/src/creator_frikshun</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>StandardOutPath</key>
  <string>/Users/allentaylor/Library/Logs/FrikShun/tunnel.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/allentaylor/Library/Logs/FrikShun/tunnel-error.log</string>
</dict>
</plist>
